# BusinessOS/04_AUTOMATIONS/systemd/businessos-memory-sync.path
[Unit]
Description=Watch canonical BusinessOS docs and sync them into Holographic memory

[Path]
PathChanged=/home/yuiop/.hermes/hermes-agent/BusinessOS/PROJECT.md
PathChanged=/home/yuiop/.hermes/hermes-agent/BusinessOS/docs/architecture.md
PathChanged=/home/yuiop/.hermes/hermes-agent/BusinessOS/docs/decisions.md
Unit=businessos-memory-sync.service

[Install]
WantedBy=default.target
