# BusinessOS/04_AUTOMATIONS/systemd/businessos-memory-sync.timer
[Unit]
Description=Periodic safety-net sync for BusinessOS memory

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h
Persistent=true
Unit=businessos-memory-sync.service

[Install]
WantedBy=timers.target
